# private_Library/private_LaunchAgents/com.piotrrojek.daily-recommendation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.piotrrojek.daily-recommendation</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/piotrrojek/Desktop/secondbrain/scripts/daily_recommendation.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/piotrrojek/Desktop/secondbrain</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>20</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/piotrrojek/Desktop/secondbrain/scripts/.daily-rec-logs/launchd.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/piotrrojek/Desktop/secondbrain/scripts/.daily-rec-logs/launchd.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/piotrrojek/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
