# private_Library/private_LaunchAgents/com.piotrrojek.news-pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.piotrrojek.news-pipeline</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>/Users/piotrrojek/Desktop/secondbrain/scripts/news_pipeline.py</string>
        <string>daily</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/piotrrojek/Desktop/secondbrain</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>0</integer>
        <key>Minute</key><integer>5</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/piotrrojek/Desktop/secondbrain/scripts/.news-logs/launchd.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/piotrrojek/Desktop/secondbrain/scripts/.news-logs/launchd.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/piotrrojek/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
