# Library/LaunchAgents/com.piotrrojek.secondbrain-chart-rotation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.piotrrojek.secondbrain-chart-rotation</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/uv</string>
    <string>run</string>
    <string>/Users/piotrrojek/Desktop/secondbrain/scripts/rotate-trading-charts.py</string>
    <string>--yes</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/piotrrojek/Desktop/secondbrain</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/piotrrojek</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>4</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/piotrrojek/Library/Logs/secondbrain-chart-rotation.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/piotrrojek/Library/Logs/secondbrain-chart-rotation.log</string>
</dict>
</plist>
